# scripts/launchd/coverage.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>{{LABEL}}</string>

  <key>ProgramArguments</key>
  <array>
    <string>{{NODE_BIN}}</string>
    <string>{{REPO_ROOT}}/scripts/run-coverage.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>{{REPO_ROOT}}</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>{{HOME}}</string>
    <key>SLACK_WEBHOOK_URL</key>
    <string>{{SLACK_WEBHOOK_URL}}</string>
  </dict>

  <key>StandardOutPath</key>
  <string>{{REPO_ROOT}}/.launchd-coverage.out.log</string>

  <key>StandardErrorPath</key>
  <string>{{REPO_ROOT}}/.launchd-coverage.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
